# com.greenheadlabs.heartbeat-v2.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.greenheadlabs.heartbeat-v2</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/dieselgoose/.openclaw/workspace/heartbeat_cron.sh</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/dieselgoose/.openclaw/workspace</string>
    
    <key>StartInterval</key>
    <integer>600</integer>
    
    <key>StandardOutPath</key>
    <string>/Users/dieselgoose/.openclaw/logs/heartbeat_v2.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/dieselgoose/.openclaw/logs/heartbeat_v2_err.log</string>
    
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
